# deploy/com.polymarket.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polymarket.server</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading/.venv/bin/python3</string>
        <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading/server.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading</string>

    <!-- Restart automatically if the process exits for any reason -->
    <key>KeepAlive</key>
    <true/>

    <!-- Wait 30 seconds before restarting to avoid spin-looping on hard crashes -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading/logs/server.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading/logs/server.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
