# deploy/com.polymarket.eod_update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polymarket.eod_update</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading/.venv/bin/python3</string>
        <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading/tools/eod_pipeline.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>20</integer>
            <key>Weekday</key>
            <integer>1</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>20</integer>
            <key>Weekday</key>
            <integer>2</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>20</integer>
            <key>Weekday</key>
            <integer>3</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>20</integer>
            <key>Weekday</key>
            <integer>4</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>13</integer>
            <key>Minute</key>
            <integer>20</integer>
            <key>Weekday</key>
            <integer>5</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading/logs/eod_update.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/jayzhuang/Desktop/Personal Agent projects/Polymarket_Tailend Trading/logs/eod_update.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
</dict>
</plist>
